# collections/ansible_collections/nvidia/bcm/collections/ansible_collections/nvidia/bcm/collections/ansible_collections/nvidia/bcm/collections/ansible_collections/nvidia/bcm/collections/ansible_collections/nvidia/bcm/collections/ansible_collections/nvidia/bcm/collections/ansible_collections/nvidia/bcm/roles/rhel_cluster/templates/bcm-agent.container.j2
[Unit]
Description=BCM Agent - NVIDIA Base Command Manager Integration
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/fabiendupont/nvidia-bcm-lite-daemon:latest
ContainerName=bcm-agent

Environment=NODE_NAME=%H
Environment=BCM_HOST={{ bcm_internal_ip }}
Environment=BCM_PORT=8081
Environment=CERT_FILE=/etc/bcm-agent/certs/cert.pem
Environment=KEY_FILE=/etc/bcm-agent/certs/cert.key
Environment=CA_FILE=/etc/bcm-agent/certs/cacert.pem
Environment=SYNC_INTERVAL=300
Environment=LABEL_PREFIX=bcm.nvidia.com
Environment=METRICS_PORT=9100
Environment=LOG_LEVEL=info

Network=host
HostName=%H
SecurityLabelDisable=true
User=0
Group=0

AddCapability=SYS_ADMIN
AddCapability=SYS_RAWIO
AddCapability=NET_ADMIN

Volume=/etc/bcm-agent/certs:/etc/bcm-agent/certs:ro,z
Volume=/dev:/dev:rw
Volume=/sys:/sys:rw
Volume=/proc:/host/proc:ro
Volume=/lib/firmware:/lib/firmware:ro
Volume=/var/log/bcm-agent:/var/log/bcm-agent:rw,z

PodmanArgs=--pid=host
PodmanArgs=--ipc=host
PodmanArgs=--privileged

HealthCmd=/bin/sh -c "test -f /var/run/bcm-agent/cm-lite-daemon.pid && kill -0 $(cat /var/run/bcm-agent/cm-lite-daemon.pid)"
HealthInterval=60s
HealthTimeout=10s
HealthRetries=3

PublishPort=9100:9100

[Service]
Restart=always
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=45
MemoryHigh=488M
MemoryMax=512M
TasksMax=128

[Install]
WantedBy=multi-user.target default.target
