# roles/bcm_agent/templates/bcm-agent.container.j2
# Quadlet configuration for BCM Agent
# Managed by Ansible - DO NOT EDIT MANUALLY
# Generated: {{ ansible_date_time.iso8601 }}

[Unit]
Description=BCM Agent - NVIDIA Base Command Manager Integration
Documentation=https://github.com/nvidia/bcm-agent
After=network-online.target
Wants=network-online.target

[Container]
# Container image
Image={{ bcm_agent_image }}

{% if bcm_agent_auto_update %}
# Automatic updates
AutoUpdate=registry
{% endif %}

# Container name
ContainerName=bcm-agent

# Environment variables
Environment=NODE_NAME=%H
Environment=BCM_HOST={{ bcm_agent_bcm_host }}
Environment=BCM_PORT={{ bcm_agent_bcm_port }}
Environment=CERT_FILE={{ bcm_agent_cert_dest }}/cert.pem
Environment=KEY_FILE={{ bcm_agent_cert_dest }}/cert.key
Environment=CA_FILE={{ bcm_agent_cert_dest }}/cacert.pem
Environment=SYNC_INTERVAL={{ bcm_agent_sync_interval }}
Environment=LABEL_PREFIX={{ bcm_agent_label_prefix }}
Environment=METRICS_PORT={{ bcm_agent_metrics_port }}
Environment=LOG_LEVEL={{ bcm_agent_log_level }}
{% if bcm_agent_debug %}
Environment=DEBUG=true
{% endif %}
{% if bcm_agent_disable_labeling | default(false) %}
Environment=DISABLE_LABELING=true
{% endif %}

# Network configuration
Network=host
HostName=%H
AddHost={{ bcm_agent_bcm_host }}:{{ bcm_agent_bcm_host_ip | default('10.141.255.254') }}

# Security and access
SecurityLabelDisable=true
User=0
Group=0

# Capabilities (for hardware access)
AddCapability=SYS_ADMIN
AddCapability=SYS_RAWIO
AddCapability=NET_ADMIN

# Volume mounts
Volume={{ bcm_agent_cert_dest }}:{{ bcm_agent_cert_dest }}:ro,z
Volume=/dev:/dev:rw
Volume=/sys:/sys:rw
Volume=/proc:/host/proc:ro
Volume=/lib/firmware:/lib/firmware:ro
Volume={{ bcm_agent_log_dir }}:{{ bcm_agent_log_dir }}:rw,z

# PID namespace
PodmanArgs=--pid=host
PodmanArgs=--ipc=host
PodmanArgs=--privileged

# Health check
HealthCmd=/bin/sh -c "test -f /var/run/bcm-agent/cm-lite-daemon.pid && kill -0 $(cat /var/run/bcm-agent/cm-lite-daemon.pid)"
HealthInterval={{ bcm_agent_health_interval }}
HealthTimeout={{ bcm_agent_health_timeout }}
HealthRetries={{ bcm_agent_health_retries }}

# Expose Prometheus metrics
PublishPort={{ bcm_agent_metrics_port }}:{{ bcm_agent_metrics_port }}

[Service]
# Restart policy
Restart={{ bcm_agent_restart_policy }}
RestartSec={{ bcm_agent_restart_sec }}
TimeoutStartSec={{ bcm_agent_timeout_start_sec }}
TimeoutStopSec={{ bcm_agent_timeout_stop_sec }}

# Resource limits
MemoryHigh={{ bcm_agent_memory_high }}
MemoryMax={{ bcm_agent_memory_max }}
TasksMax={{ bcm_agent_tasks_max }}

[Install]
WantedBy=multi-user.target default.target
